# lib/deploy/resources/services/unicorn/run
#!/bin/sh
exec envdir ./env setuidgid hyperarchy ./unicorn.sh 2>&1
