# lib/deploy/resources/services/unicorn/run
#!/bin/sh
exec 2>&1
exec envdir ./env setuidgid hyperarchy ./unicorn.sh
